# flyspray/fly_dhb/supervisor/fly_dhb.conf
[program:fly_dhb]
command=/opt/fly-dhb/run.py
directory=/opt/fly-dhb/
autostart=true
autorestart=true
redirect_stderr=true
environment=PYTHONUNBUFFERED="1"