# flyspray/fly_export/supervisor/fly_export.conf
[program:fly_export]
command=/opt/fly-export/run.py
directory=/opt/fly-export/
autostart=true
autorestart=true
redirect_stderr=true
environment=PYTHONUNBUFFERED="1"