# victoriametrics/tools/vm.conf
[program:vm]
command=/opt/vm/run.sh
directory=/opt/vm/
user=pi
autostart=true
autorestart=true
redirect_stderr=true